# deploy/supervisor/ams-worker.conf
; =============================================================================
; Supervisor — Laravel queue workers for AMS
; Copy to /etc/supervisor/conf.d/ams-worker.conf then:
;   supervisorctl reread && supervisorctl update && supervisorctl start ams-worker:*
; =============================================================================

[program:ams-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ams/backend/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --backoff=10
directory=/var/www/ams/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4                          ; scale workers for 50k+ users
redirect_stderr=true
stdout_logfile=/var/log/ams/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=3600

; Dedicated worker for the (higher-volume) notifications queue.
[program:ams-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ams/backend/artisan queue:work redis --queue=notifications,default --sleep=2 --tries=3 --max-time=3600
directory=/var/www/ams/backend
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/ams/notifications.log
